# ops/dogv-embed.service
[Unit]
Description=DOGV AI embedding server (llama.cpp bge-m3, GPU 1)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=exec
User=ai-server
WorkingDirectory=/home/ai-server/dogv-ai
Environment="CUDA_VISIBLE_DEVICES=1"
# --verbosity 0: default per-request logging destabilized the server under
# sustained load (crashed ~every 40min, log ballooned to 61GB). -ub 4096 holds a
# full embed request (8 chunks x ~500 tok) while keeping GPU memory ~1GB so it
# coexists with the chat server.
ExecStart=/home/ai-server/llama.cpp/build/bin/llama-server \
  --model /home/ai-server/models/bge-m3/bge-m3-f16.gguf \
  --alias bge-m3 \
  --host 127.0.0.1 --port 8001 \
  --embeddings \
  --pooling cls \
  --ctx-size 8192 -b 8192 -ub 4096 \
  --n-gpu-layers 99 \
  --embd-normalize 2 \
  --verbosity 0
ExecStartPost=/bin/bash -c 'for _ in $(seq 1 60); do curl -fsS -m 5 http://127.0.0.1:8001/health >/dev/null 2>&1 && exit 0; sleep 1; done; echo "embed server not healthy after 60s"; exit 1'
TimeoutStartSec=90
Restart=on-failure
RestartSec=10
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
