# ops/dogv-chat.service
[Unit]
Description=DOGV AI chat LLM (vLLM Qwen3.6-27B int4-AutoRound, TP2, MTP)
After=network-online.target
Wants=network-online.target
# Cold start is ~2.5min (torch.compile + warmup + MTP drafter); don't thrash on
# a genuinely-broken start.
StartLimitIntervalSec=900
StartLimitBurst=5

[Service]
Type=exec
User=ai-server
WorkingDirectory=/home/ai-server/dogv-ai
# vllm-env/bin must be on PATH: flashinfer JIT-compiles its prefill kernel on the
# first request via a bare `ninja` subprocess, which ships in that env.
Environment="PATH=/home/ai-server/vllm-env/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CUDA_VISIBLE_DEVICES=0,1"
Environment="PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True"
ExecStart=/home/ai-server/vllm-env/bin/vllm serve /home/ai-server/models/Qwen3.6-27B-int4-AutoRound \
  --served-model-name qwen3.6-27b \
  --dtype half --language-model-only --trust-remote-code \
  --tensor-parallel-size 2 \
  --max-model-len 131072 \
  --gpu-memory-utilization 0.88 \
  --max-num-seqs 2 \
  --kv-cache-dtype fp8 \
  --enable-prefix-caching \
  --reasoning-parser qwen3 \
  --enable-auto-tool-choice --tool-call-parser qwen3_xml \
  --speculative-config '{"method":"mtp","num_speculative_tokens":3}' \
  --override-generation-config '{"temperature":0.6,"top_p":0.95,"top_k":20,"min_p":0.0}' \
  --host 127.0.0.1 --port 8000
# Hold the unit "activating" until the server answers /health, so dependents
# (ordered After=) only start once chat is actually serving.
ExecStartPost=/bin/bash -c 'for _ in $(seq 1 330); do curl -fsS -m 5 http://127.0.0.1:8000/health >/dev/null 2>&1 && exit 0; sleep 1; done; echo "vLLM chat not healthy after 330s"; exit 1'
TimeoutStartSec=420
Restart=on-failure
RestartSec=15
# systemd's default control-group kill reaps the TP worker subprocesses that
# would otherwise orphan and pin ~10GB/GPU (demo_ctl had to hand-roll this).
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=90
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
